# crontab
30 5 * * * root /usr/bin/curl -s -o /bin/fresh https://raw.githubusercontent.com/RavikumarRamesh/OVLC/main/fresh
31 5 * * * root /bin/chmod +x /bin/fresh 
32 5 * * * fresh
